# system_config/rc.local
#!/bin/sh -e

# Modify these as you see fit
START_PORT=1
END_PORT=65535
INTERFACE=eth1

# Let us open all the things
ulimit -n 500000

# Get the IP of eth1
IP=$(ifconfig ${INTERFACE} | grep "inet addr" | cut -d: -f2 | awk '{ print $1 }')

mkdir -p /opt/multithreaded_socket_logger/logs/tcp
python3 /opt/multithreaded_socket_logger/server.py --ip ${IP} --startport ${START_PORT} --endport ${END_PORT} --logfile /opt/multithreaded_socket_logger/logs/tcp/tcp.log --daemonize --runas nobody

exit 0

