# system_config/rc.local
#! /bin/bash

# Modify these as you see fit
START_PORT=6000
END_PORT=6100
INTERFACE=eth1

# Get the IP of eth1
IP=$(ifconfig ${INTERFACE} | grep "inet addr" | cut -d: -f2 | awk '{ print $1 }')

mkdir -p /opt/multithreaded_socket_logger/logs/tcp
python /opt/multithreaded_socket_logger/server.py --ip ${IP} --startport ${START_PORT} --endport ${END_PORT} --logfile /opt/multithreaded_socket_logger/logs/tcp/tcp.log --daemonize --runas nobody

exit 0

